# fpm-nginx/supervosprd-php-fpm.conf
[program:php-fpm]
process_name=%(program_name)s_%(process_num)02d
command=php-fpm -F
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
